# server/infra/upstart/node-stop.conf
#!upstart
description "Node stop"

start on stopping node

task

script
    stop node-instance PORT=8081
    stop node-instance PORT=8082
    stop node-instance PORT=8083
end script