# server/infra/upstart/node-start.conf
#!upstart
description "Node start"

start on starting node

task

script
    start node-instance PORT=8081
    #start node-instance PORT=8082
    #start node-instance PORT=8083
end script